# multiplexd.socket
[Unit]
Description=Multiplex multiple services onto a single port socket

[Socket]
ListenStream=443

[Install]
WantedBy=sockets.target

